# backend/workbench.service
[Unit]
Description=Workbench Backend
After=network.target

[Service]
Type=simple
User=vault
WorkingDirectory=/home/vault/workbench/backend
ExecStart=/home/vault/workbench/backend/.venv/bin/uvicorn src.main:app --host 0.0.0.0 --port 8765
Restart=on-failure
RestartSec=10
Environment="PATH=/home/vault/workbench/backend/.venv/bin:/usr/bin"

[Install]
WantedBy=multi-user.target
