# deploy/systemd/opencode-serve.service
[Unit]
Description=OpenCode Headless Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=active-workbench
Group=active-workbench
WorkingDirectory=/opt/active-workbench
EnvironmentFile=/etc/active-workbench/opencode.env
Environment=PATH=/home/active-workbench/.local/bin:/usr/local/bin:/usr/bin
Environment=OPENCODE_SERVER_HOSTNAME=0.0.0.0
ExecStart=/usr/bin/env opencode serve --hostname ${OPENCODE_SERVER_HOSTNAME} --port 4096
Restart=always
RestartSec=2
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/var/lib/active-workbench /opt/active-workbench /home/active-workbench

[Install]
WantedBy=multi-user.target
