# deploy/systemd/vector.service.d/override.conf.example
[Service]
EnvironmentFile=/etc/active-workbench/vector.env

